# root/etc/s6-overlay/s6-rc.d/qbittorrent-config/up
#!/usr/bin/execlineb -P

with-contenv
importas alternative WEBUI_ALTERNATIVE

define confdir /container/profile/qBittorrent/config
define conf ${confdir}/qBittorrent.conf

foreground {
  if { eltest ! -d $confdir }
  mkdir -p $confdir
}

if { eltest ! -f $conf }
redirfd -x 1 $conf
echo
"
[LegalNotice]
Accepted=true

[Preferences]
WebUI\\RootFolder=/opt/webui
WebUI\\AlternativeUIEnabled=${alternative}
"
