# root/etc/s6-overlay/s6-rc.d/qbittorrent/run
#!/usr/bin/execlineb -P

with-contenv
importas port WEBUI_PORT

qbittorrent-nox --webui-port=${port} --profile=/container/profile
